# deploy/ubuntu/stock-signal-telegram-bot.service
# Шаблон systemd для Telegram-бота (Ubuntu VPS).
# 1) Скопируйте проект, например: /opt/stock_signal_analyzer
# 2) Создайте venv, pip install -r requirements.txt (+ опционально requirements-tbank.txt)
# 3) Настройте .env (TELEGRAM_BOT_TOKEN обязателен; см. .env.example)
# 4) Подставьте пути ниже и: sudo cp ... /etc/systemd/system/
#    sudo systemctl daemon-reload && sudo systemctl enable --now stock-signal-telegram-bot
#
# Пользователь: отдельный аккаунт лучше, чем root (создайте: sudo adduser --system --group --home /opt/stock_signal_analyzer stockbot)

[Unit]
Description=Stock Signal Analyzer — Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stockbot
Group=stockbot
WorkingDirectory=/opt/stock_signal_analyzer
Environment=PYTHONUNBUFFERED=1
# Секреты: EnvironmentFile для systemd; приложение также читает .env через stenv.load_project_env() (stenv.py + python-dotenv)
EnvironmentFile=-/opt/stock_signal_analyzer/.env
ExecStart=/opt/stock_signal_analyzer/.venv/bin/python /opt/stock_signal_analyzer/telegram_bot.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

# Базовое ужесточение (при необходимости отключите строки, если что-то ломается)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
