# deploy/stock-signal-bot.service
[Unit]
Description=Stock Signal Analyzer - Telegram Bot
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKING_DIR%
Environment="PATH=%WORKING_DIR%/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=%WORKING_DIR%/.env
ExecStart=%WORKING_DIR%/venv/bin/python %WORKING_DIR%/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=append:%WORKING_DIR%/logs/bot.log
StandardError=append:%WORKING_DIR%/logs/bot.error.log

# Ограничения ресурсов
MemoryLimit=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
