# deploy/stock-signal-tracker.timer
[Unit]
Description=Stock Signal Analyzer - Outcome Tracker Timer
Requires=stock-signal-tracker.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=stock-signal-tracker.service

[Install]
WantedBy=timers.target
