# deploy/stock-signal-bot-simple.service
[Unit]
Description=Stock Signal Analyzer Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/stock_signal_analyzer
Environment="PATH=/root/stock_signal_analyzer/venv/bin"
ExecStart=/root/stock_signal_analyzer/venv/bin/python telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
