# deploy/stock-signal-tracker.service
[Unit]
Description=Stock Signal Analyzer - Outcome Tracker
After=network.target

[Service]
Type=oneshot
User=%USER%
WorkingDirectory=%WORKING_DIR%
Environment="PATH=%WORKING_DIR%/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=%WORKING_DIR%/.env
ExecStart=%WORKING_DIR%/venv/bin/python -m stock_signal_analyzer.outcome_tracker
StandardOutput=append:%WORKING_DIR%/logs/tracker.log
StandardError=append:%WORKING_DIR%/logs/tracker.error.log
